# backend/scripts/com.somm.extraction-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.somm.extraction-enrich</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/danielguerrero/dev/wine_app/backend/scripts/run_extraction_launchd.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/danielguerrero/dev/wine_app/backend</string>

    <!-- Weekly, Sunday 03:00 local time (after weekly-scrape at 02:00 CT).
         If asleep at the trigger, launchd runs it on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/danielguerrero/Library/Logs/somm-extraction.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielguerrero/Library/Logs/somm-extraction.err.log</string>

    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
