# backend/scripts/com.somm.twin-liquors.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.somm.twin-liquors</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/danielguerrero/dev/wine_app/backend/scripts/run_twin_liquors_launchd.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/danielguerrero/dev/wine_app/backend</string>

    <!-- Weekly Sun 04:00 CT — after the GitHub weekly-scrape (Sun 02:00 CT)
         and the local extraction LaunchAgent (Sun 03:00 CT), so this
         doesn't fight either for CPU/DB. Twin Liquors runs on the mini
         because GitHub datacenter IPs are Cloudflare-1015-blocked. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>4</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/danielguerrero/Library/Logs/somm-twin-liquors.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/danielguerrero/Library/Logs/somm-twin-liquors.err.log</string>

    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
